# runtime/elders/systemd/ttyd-elder.service.template
[Unit]
Description=ttyd readonly web view for elder-{{N}} tmux session
Documentation=https://github.com/tsl0922/ttyd
After=default.target

[Service]
Type=simple
Environment=TZ=America/New_York
ExecStart=%h/bin/ttyd -R -p {{PORT}} -t fontSize=11 -t cursorStyle=bar tmux attach-session -t elder-{{N}}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
